# opentrace.service
[Unit]
Description=OpenTrace
After=network.target

[Service]
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/opentrace
ExecStart=/home/YOUR_USER/opentrace/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
